# release/linux/trezord.service
[Unit]
Description=TREZOR Bridge

[Service]
Type=simple
ExecStart=/usr/bin/trezord
