# release/linux/trezord.service
[Unit]
Description=TREZOR Bridge
After=multi-user.target

[Service]
Type=forking
ExecStart=/usr/bin/trezord

[Install]
WantedBy=multi-user.target
